# move_me/etc/systemd/system/alfred.service
[Unit]
Description=ALFRED primary on bat0
After=network-online.target
Wants=network-online.target

# warte kurz, bis bat0 existiert (z. B. durch dein batmesh.sh)
# wenn du einen eigenen Service für das Mesh hast, trag ihn bei After= ein
[Service]
Type=simple
ExecStartPre=/bin/sh -c 'for i in $(seq 1 20); do ip link show bat0 >/dev/null 2>&1 && exit 0; sleep 1; done; echo "bat0 not ready"; exit 1'
ExecStart=/usr/sbin/alfred -i bat0 -m
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
