# move_me/etc/systemd/system/mesh-startup.service
# /etc/systemd/system/mesh-startup.service
[Unit]
Description=Mesh Network and Monitor (startup_sequence.sh wrapper)
After=local-fs.target systemd-modules-load.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Definierter PATH, damit ip/iw/batctl gefunden wird
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Falls /home separat gemountet ist
WorkingDirectory=/home/natak
ExecStart=/home/natak/mesh/startup_sequence.sh
# 60s auf Mesh-Join warten (reduziert von 180s)
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
